# docker/s6/services/lazyraster.svc/run
#!/bin/sh

echo "Launching Lazyraster service..."
/lazyraster/lazyraster
